# agent/sreoob-agent.service
[Unit]
Description=SREoob Monitoring Agent
Documentation=https://github.com/your-org/sreoob
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
User=sreoob
Group=sreoob
WorkingDirectory=/opt/sreoob-agent
ExecStart=/opt/sreoob-agent/sreoob-agent
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Essential Environment Variables
Environment=MASTER_FQDN=https://your-sreoob-server.com
Environment=AGENT_API_KEY=your_64_character_agent_api_key_here
Environment=AGENT_ID=agent-prod-01

# Optional Environment Variables (with defaults)
Environment=AGENT_SERVER_PORT=5227
Environment=AGENT_PORT=8081
Environment=USE_WEBSOCKET=true
Environment=LOG_LEVEL=info

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/sreoob-agent
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource Limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sreoob-agent

[Install]
WantedBy=multi-user.target 